# arch-setup/etc/systemd/system/power-profiles-daemon.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/lib/power-profiles-daemon --block-action=amdgpu_panel_power
