# arch-setup/etc/systemd/system/dlfcd.service
[Unit]
Description=Dell laptop thermal monitor and control fan

[Service]
ExecStart=/usr/local/bin/dlfcd

[Install]
WantedBy=multi-user.target
