# arch-setup/etc/systemd/system/iwd.service.d/override.conf
[Service]
ExecStartPre=/usr/bin/sleep 2
